# linmon.conf
# ============================================================================
# LinMon Configuration File - Production Security Configuration
# ============================================================================
# This is a balanced production configuration that enables:
# - Core monitoring (processes, network)
# - Binary hashing and package verification
# - Low-noise security detections (credential theft, rootkits, fileless malware)
# - Reasonable filtering (private networks, thread noise)
#
# Optimized for: Security monitoring with acceptable performance overhead
# Target systems: Production servers, bastion hosts, critical infrastructure

# ============================================================================
# LOGGING CONFIGURATION
# ============================================================================

# Path to JSON event log file
# Events are logged in JSON Lines format (one JSON object per line)
log_file = /var/log/linmon/events.json

# Enable logging ALL events to syslog (in addition to JSON file)
# true  = All events are logged to both JSON file AND syslog/journald
# false = Only daemon lifecycle events go to syslog (default)
# Note: Daemon lifecycle events are ALWAYS logged to syslog for tamper detection
# Recommendation: false (use remote syslog forwarding for tamper resistance instead)
log_to_syslog = false

# ============================================================================
# BUILT-IN LOG ROTATION
# ============================================================================

# Enable built-in log rotation
# true  = LinMon rotates logs automatically
# false = Use external logrotate (linmond.logrotate)
log_rotate = true

# Maximum log file size before rotation
# Supports suffixes: K (KB), M (MB), G (GB)
# Default: 100M
log_rotate_size = 100M

# Number of rotated log files to keep
# Rotation pattern: events.json -> events.json.1 -> events.json.2 -> ...
# Default: 10
log_rotate_count = 10

# Verbosity level (0=quiet, 1=normal, 2=verbose)
verbosity = 1

# ============================================================================
# MONITORING CONFIGURATION
# ============================================================================

# Monitor process execution and termination
# Logs: process_exec, process_exit events
# This is the core of LinMon - always enabled
monitor_processes = true

# Monitor process exit events
# true  = Log both exec and exit events (more complete audit trail)
# false = Only log exec events (reduces noise from short-lived processes)
# Recommendation: false for production (exit events rarely needed)
monitor_process_exit = false

# Monitor file operations (create, modify, delete)
# Warning: Can generate VERY high event volume on busy systems
# Logs: file_create, file_modify, file_delete events
# Recommendation: false (enable only on specific systems with ignore_file_paths)
monitor_files = false

# Monitor TCP connections (connect/accept)
# Logs: net_connect_tcp, net_accept_tcp events
# Supports both IPv4 and IPv6
# Recommendation: true (critical for lateral movement detection)
monitor_tcp = true

# Monitor UDP traffic (sendmsg)
# Warning: UDP can be VERY noisy (DNS, NTP, DHCP, mDNS, etc.)
# Logs: net_send_udp events
# Recommendation: false (enable only if you need DNS/NTP monitoring)
monitor_udp = false

# Monitor vsock (Virtual Socket) connections
# vsock is used for VM-to-host and container-to-host communication
# Logs: net_vsock_connect events with CID (Context ID)
# Use case: Detect container escape attempts (MITRE ATT&CK T1611)
# Recommendation: false (enable only on VM/container hosts)
monitor_vsock = false

# ============================================================================
# FILTERING CONFIGURATION
# ============================================================================

# UID filtering - only log events from users with UID >= min_uid
# 0    = Log all users including root (RECOMMENDED for security)
# 1000 = Only log human users (reduces system service noise)
# Recommendation: 0 (critical to detect root compromise)
min_uid = 0

# Maximum UID to monitor (0 = no limit)
max_uid = 0

# Require controlling TTY (terminal session)
# true  = Only log interactive terminal sessions (ssh, local terminal)
# false = Log all processes including GUI, systemd services
# Recommendation: false (Sysmon-like comprehensive monitoring)
require_tty = false

# Ignore threads (only log main processes)
# true  = Filter out threads (reduces browser/GUI noise significantly)
# false = Log everything (complete audit trail)
# Recommendation: true (reduces Firefox/Chrome thread pool noise by 90%+)
ignore_threads = true

# ============================================================================
# PROCESS FILTERING
# ============================================================================

# Process name blacklist (comma-separated, exact match on comm field)
# Use this to reduce noise from known-safe chatty processes
# Examples: systemd,dbus-daemon,chrome,firefox
# Recommendation: Leave empty, use SIEM-side filtering instead
ignore_processes =

# Process name whitelist (only log these if set)
# Leave empty to log all processes
# Use case: Focused monitoring of specific processes (sudo, ssh, etc.)
only_processes =

# ============================================================================
# NETWORK & FILE PATH FILTERING
# ============================================================================

# Network CIDR filtering (comma-separated IPv4 CIDR blocks to ignore)
# Events TO these destinations are dropped in kernel (high performance)
# Note: Only filters destination address, not source
# Recommendation: Filter private networks to focus on external connections
ignore_networks = 127.0.0.0/8,10.0.0.0/8,172.16.0.0/12,192.168.0.0/16

# File path filtering (comma-separated path prefixes to ignore)
# Only relevant if monitor_files = true
# Examples: /tmp/,/proc/,/sys/,/dev/,.cache/
ignore_file_paths =

# ============================================================================
# ENRICHMENT & PRIVACY
# ============================================================================

# Capture full command-line arguments
# true  = Capture argv[] for all process_exec events (RECOMMENDED)
# false = Only capture process name
# Note: Combined with redact_sensitive, this is safe
capture_cmdline = true

# Redact sensitive patterns in command lines
# Automatically removes passwords, tokens, API keys from cmdline
# Patterns: password=, token=, api_key=, -p, --password, etc.
# Recommendation: true (critical privacy protection)
redact_sensitive = true

# Resolve UID to username (adds "username" field to events)
# Uses cached /etc/passwd lookups (256 entry LRU cache)
# Performance: negligible overhead (~1μs per lookup after cache hit)
# Recommendation: true (improves SIEM readability)
resolve_usernames = true

# Calculate SHA256 hash of executed binaries
# Adds "sha256" field to process_exec events
# Uses persistent disk cache (survives daemon restarts)
# Performance: ~1-5ms first hash, <1μs cached lookups
# Recommendation: true (critical for threat hunting and IOC matching)
hash_binaries = true

# Verify if executed binaries belong to system packages
# Adds "package" field (package name or null) to process_exec events
# Adds "pkg_modified" field (true if binary modified since install)
# Uses persistent disk cache with 24-hour TTL (v1.3.1+)
# Performance: ~10-50ms first lookup, <1μs cached lookups
# Recommendation: true (detect unpackaged/modified binaries)
verify_packages = true

# ============================================================================
# CACHE CONFIGURATION
# ============================================================================
# LinMon caches SHA256 hashes and package lookups to disk for persistence
# across daemon restarts. Caches are validated using inode+mtime+TTL.

# SHA256 hash cache file path (used when hash_binaries = true)
# Default: /var/cache/linmon/hashes.cache
# hash_cache_file = /var/cache/linmon/hashes.cache

# Maximum entries in SHA256 hash cache
# Range: 100-1000000, Default: 10000
# hash_cache_size = 10000

# Package cache file path (used when verify_packages = true)
# Default: /var/cache/linmon/packages.cache
# pkg_cache_file = /var/cache/linmon/packages.cache

# Maximum entries in package cache
# Range: 100-1000000, Default: 10000
# pkg_cache_size = 10000

# Periodic cache save interval (minutes)
# 0  = Only save at shutdown (risk of data loss on crash)
# >0 = Save every N minutes AND at shutdown/reload
# Range: 0-60, Default: 5
cache_save_interval = 5

# ============================================================================
# TAMPER DETECTION & INTEGRITY MONITORING
# ============================================================================

# Periodic integrity checkpoint interval (minutes)
# LinMon logs checkpoints to syslog/journald with:
#   - Event sequence number (detect deleted events via gaps)
#   - LinMon binary SHA256 (detect daemon replacement)
#   - Config file SHA256 (detect unauthorized config changes)
# Range: 0-1440 minutes (0 = disabled)
# Recommendation: 30 (checkpoint every 30 minutes)
checkpoint_interval = 30

# ============================================================================
# SECURITY MONITORING - LOW NOISE (RECOMMENDED - ENABLED)
# ============================================================================
# These detections have very low false positive rates and catch critical
# attack techniques. Minimal performance overhead.

# Monitor credential/auth file reads - T1003.008, T1552.004, T1098.004
# Detects: Reads of /etc/shadow, /etc/sudoers, /etc/ssh/*, /etc/pam.d/*
#          Reads of ~/.ssh/id_* (private keys), ~/.ssh/authorized_keys, ~/.ssh/config
# Smart filtering: Whitelists legitimate auth processes (login, sshd, sudo, etc.)
# Performance: negligible (only fires on suspicious reads)
# Known limitation: SSH keys only detected if /.ssh/ within first 22 chars (covers 95%+ systems)
# Recommendation: true (CRITICAL - catches credential theft and SSH key attacks)
monitor_cred_read = true

# Monitor /etc/ld.so.preload writes - T1574.006 LD_PRELOAD Hijacking
# Detects: Writes to /etc/ld.so.preload (library injection rootkit)
# Performance: negligible (file is almost never written legitimately)
# Recommendation: true (CRITICAL - almost zero noise, catches rootkits)
monitor_ldpreload = true

# Monitor kernel module loading - T1547.006 Rootkit/Driver detection
# Detects: init_module(), finit_module() syscalls
# Use case: Detect kernel rootkits, malicious drivers
# Performance: negligible (modules loaded rarely on production systems)
# Recommendation: true (HIGH VALUE - alert on any module load)
monitor_modules = true

# Monitor memfd_create() - T1620 Fileless Malware detection
# Detects: Anonymous executable memory regions
# Use case: Fileless malware, in-memory ELF execution
# Performance: low overhead (only fires on memfd creation)
# Note: Some containers/JIT compilers use memfd legitimately
# Recommendation: true (correlate with execveat for full fileless chain)
monitor_memfd = true

# Monitor execveat() syscall - T1620 Fileless Execution detection
# Detects: fd-based execution (executing from file descriptor)
# Use case: Fileless malware using memfd_create() + execveat() pattern
# Performance: low overhead (rare syscall)
# Recommendation: true (completes fileless malware detection)
monitor_execveat = true

# Monitor bpf() syscall - T1014 eBPF Rootkit detection
# Detects: BPF program loading, map creation
# Use case: eBPF-based rootkits, packet manipulation
# Performance: low overhead
# Note: Legitimate tools use eBPF (bpftool, tc, ip, LinMon itself)
# Recommendation: true (whitelist known tools in SIEM)
monitor_bpf = true

# ============================================================================
# SECURITY MONITORING - HIGHER NOISE (OPTIONAL - DISABLED)
# ============================================================================
# These detections can have false positives from legitimate tools.
# Enable based on threat model and accept tuning overhead.

# Monitor ptrace() syscall - T1055 Process Injection detection
# Detects: PTRACE_ATTACH, PTRACE_POKETEXT, PTRACE_POKEDATA
# Use case: Process injection, debugger attachment
# Note: gdb, strace, ltrace trigger this legitimately
# Recommendation: false (enable on systems without debuggers)
monitor_ptrace = false

# Monitor bind() syscall - T1571 Bind Shell / C2 Server detection
# Detects: Processes binding to network ports
# Use case: Bind shells, unauthorized services
# Note: Web servers, databases, legitimate services trigger this
# Recommendation: false (enable with SIEM-side process name correlation)
monitor_bind = false

# Monitor unshare() syscall - T1611 Container Escape detection
# Detects: Namespace manipulation
# Use case: Container escapes, user namespace privilege escalation
# Note: Docker/Podman use unshare legitimately
# Recommendation: false (enable on container hosts with known baseline)
monitor_unshare = false

# Monitor persistence mechanisms - T1053, T1547 Persistence detection
# Detects: Writes to cron, systemd services, shell profiles, init scripts, autostart
# Locations: /etc/cron.d/*, /var/spool/cron/*, /etc/systemd/system/*,
#            ~/.bashrc, ~/.profile, ~/.zshrc, ~/.config/autostart/*
# Use case: Detect attackers establishing persistence after compromise
# Note: Legitimate admin tasks write to these locations
# Known limitation: Shell profiles only detected if pattern within first 38 chars (covers 95%+ systems)
# Recommendation: false (enable on production servers with change management)
monitor_persistence = false

# Monitor SUID/SGID manipulation - T1548.001 Privilege Escalation detection
# Detects: chmod operations setting SUID (04000) or SGID (02000) bits
# Use case: Detect privilege escalation setup (chmod +s for persistent root access)
# Note: Package managers may set SUID on some binaries during install
# Recommendation: false (enable on high-security systems, alert on unexpected paths)
monitor_suid = false

# Monitor credential file WRITES - T1098.001 Account Manipulation detection
# Detects: Writes to /etc/passwd, /etc/shadow, /etc/sudoers, ~/.ssh/authorized_keys
# Use case: Account creation, password changes, SSH backdoors
# Recommendation: true (critical detection, very low noise)
monitor_cred_write = true

# Monitor log file tampering - T1070.001 Log Clearing detection
# Detects: Truncation or deletion of /var/log/* files
# Use case: Detect attackers covering their tracks
# Recommendation: true (critical detection, whitelists legitimate log managers)
monitor_log_tamper = true

# ============================================================================
# NOTES
# ============================================================================
# This configuration is optimized for:
# - Security visibility (credential theft, rootkits, lateral movement)
# - Performance (low overhead, <2% CPU typical)
# - Low noise (minimal false positives, ~100-500 events/day typical)
#
# Expected event volume (typical production server):
# - Process exec: 50-200/day (depending on workload)
# - Network TCP: 100-500/day (external connections only, private filtered)
# - Security events: 0-10/day (should be investigated)
# - Total: ~200-1000 events/day (~14KB-70KB/day compressed)
#
# Performance impact:
# - CPU: 0.5-2% (typical), 5-8% (high exec rate like compile servers)
# - Memory: ~50MB daemon + BPF maps
# - Disk I/O: ~10-50MB/day logs (highly variable)
#
# Recommended tuning for specific environments:
# - Desktop/laptop: Set ignore_threads=true, min_uid=1000
# - Database server: Set monitor_files=false
# - Container host: Set monitor_vsock=true, monitor_unshare=true
# - Bastion host: Set verify_packages=true, alert on package=null
# - Compile server: Increase rate limits or use ignore_processes
#
# For SIEM integration examples, see:
# - extras/filebeat/filebeat-ecs.yml (ECS format for Elasticsearch)
# - extras/vector/vector.toml (ClickHouse for analytics)
# - extras/rsyslog-remote.conf (Tamper-resistant remote logging)
